# webapp.service
[Unit]
Description=Flask Web Application
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Sean-Anthony-Mael
Environment="PATH=/home/ubuntu/Sean-Anthony-Mael/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ubuntu/Sean-Anthony-Mael/.venv/bin/python3 /home/ubuntu/Sean-Anthony-Mael/src/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webapp

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

